# .claude/VoiceServerV1/com.pai.voiceserver-v1.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pai.voiceserver-v1</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BUN__</string>
        <string>run</string>
        <string>__HOME__/.claude/VoiceServerV1/server.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/.claude/VoiceServerV1</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8888</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__BUN_DIR__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude/VoiceServerV1/voice-server.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/VoiceServerV1/voice-server-error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
